# tcp_forward.service
[Unit]
Description=TCP Port Forwarding Service
After=network.target

[Service]
ExecStart=/root/tcp_forward.sh
Restart=always
RestartSec=3
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
